# daemon/rog-daemon.service
[Unit]
Description = ROG control daemon
After = systemd-udevd.service

[Service]
ExecStart = /usr/bin/rog-daemon
Type = simple
Restart = always

[Install]
WantedBy = multi-user.target
